# systemd-timer/bunsen-pepperflash.timer
[Unit]
Description=Run Adobe pepperflash upgrade check daily.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
